# gmenuproxy/cutefish-gmenuproxy.service.in
[Unit]
Description=Proxies GTK DBus menus to a Cutefish readable format
PartOf=graphical-session.target

[Service]
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/cutefish-gmenuproxy
Restart=on-failure
Type=simple
Slice=background.slice
TimeoutSec=5sec
